# deploy/ess_guardian_runtime/systemd/etc/systemd/system/ess-vision.service
[Unit]
Description=ESS Vision Demo Service
ConditionPathExists=/dev/cam_rgb
After=network.target ess-camera-setup.service ess-thermal.service
Wants=ess-camera-setup.service ess-thermal.service
PartOf=ess.target

[Service]
Type=simple
ExecStart=/opt/ess-guardian/current/bin/ess-vision-demo
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=ess.target
