# deploy/ess_guardian_runtime/systemd/etc/systemd/system/ess-thermal.service

[Unit]
Description=ESS Thermal Daemon (MLX90640 -> status json)
After=network-online.target systemd-udevd.service
Wants=network-online.target systemd-udevd.service
PartOf=ess.target

[Service]
Type=simple
Environment=ESS_THERMAL_I2C_BUS=9
ExecStart=/opt/ess-guardian/current/thermal/venv/bin/python /opt/ess-guardian/current/thermal/ess_thermal_daemon.py
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=ess.target
